# packaging/ubuntu-vivid64/log-reader/zipnish.service
[Unit]
Description=Zipnish log reader service

[Service]
ExecStart=/usr/share/zipnish/log-reader/venv/bin/python /usr/share/zipnish/log-reader/app.py