# ui/redhat/zipnish-ui.service
[Unit]
Description=Zipnish user interface.

[Service]
WorkingDirectory=/opt/zipnish-ui
ExecStart=/opt/zipnish-ui/manage.py runserver
