# ui/redhat/zipnish-ui.service
[Unit]
Description=Zipnish user interface.

[Service]
WorkingDirectory=/opt/zipnish/ui
ExecStart=/usr/bin/python /opt/zipnish/ui/app.py
