# packaging/ubuntu-vivid64/debian/zipkin-ui.service
[Unit]
Description=Zipnish user interface service

[Service]
ExecStart=/usr/share/zipnish/ui/venv/bin/python /usr/share/zipnish/ui/app.py