# provisioning/templates/backend/siege-varnish.service.j2
[Unit]
Description=Run siege on varnish {{ varnish_address }}:{{ varnish_port }}/process-order

[Service]
ExecStart=/usr/bin/siege -c 1 {{ varnish_address }}:{{ varnish_port }}/process-order