# log-reader/redhat/log-reader.service
[Unit]
Description=Log reader service (varnish 4.0)

[Service]
WorkingDirectory=/opt/zipnish
ExecStart=/opt/zipnish/app.py
